# src/meteod.service
[Unit]
Description=Meteo Daemon
After=network.target

[Service]
Type=simple
User=meteo
WorkingDirectory=/home/meteo/meteo
ExecStart=/home/meteo/meteo/meteod 
RestartSec=10
Restart=always
StandardOutput=null
StandardError=file:/var/log/meteod.err

[Install]
WantedBy=multi-user.target
